# deploy/systemd/webex-config-pull-worker.service
[Unit]
Description=Webex config pull worker
Documentation=man:systemd.service(5)
Wants=network-online.target
After=network-online.target
ConditionPathExists=/sys/fs/cgroup/cgroup.controllers

[Service]
Type=simple
User=webex-config-deploy
Group=webex-config-pull
SupplementaryGroups=webex-config-deploy
WorkingDirectory=/
ExecStart=/usr/bin/node /opt/webex-generic-account-bot/code/scripts/config-pull-worker.mjs
Restart=on-failure
RestartSec=5s
TimeoutStopSec=15s
KillMode=control-group
ReadOnlyPaths=/run/webex-config-deploy
ReadWritePaths=/run/webex-config-deploy/config-pull-worker.lock
ReadWritePaths=/run/webex-config-deploy/deploy-config.lock
ReadWritePaths=/run/webex-config-pull
ReadWritePaths=/var/lib/webex-generic-account-bot/config-actions
ReadWritePaths=/var/lib/webex-generic-account-bot/config-prepare-checkout
ReadWritePaths=/var/lib/webex-generic-account-bot/config-staging
ReadOnlyPaths=-/var/lib/webex-generic-account-bot/rendered
UMask=0077

NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
PrivateDevices=true
ProtectClock=true
ProtectControlGroups=true
ProtectHostname=true
ProtectKernelLogs=true
ProtectKernelModules=true
ProtectKernelTunables=true
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
RestrictNamespaces=true
RestrictRealtime=true
RestrictSUIDSGID=true
LockPersonality=true
CapabilityBoundingSet=
AmbientCapabilities=
SystemCallArchitectures=native
TasksMax=256
LimitNOFILE=1024
LimitNPROC=256
MemoryMax=1G
StandardInput=null

[Install]
WantedBy=multi-user.target
