# deploy/systemd/webex-codex-launcher@.service
[Unit]
Description=Webex Codex launcher connection
Documentation=man:systemd.service(5)
ConditionPathExists=/sys/fs/cgroup/cgroup.controllers
Requires=webex-codex-launcher.socket
After=webex-codex-launcher.socket
CollectMode=inactive-or-failed

[Service]
Type=exec
User=root
Group=root
Slice=system.slice
SupplementaryGroups=webex-codex-input webex-codex-launch
ExecStart=/opt/webex-generic-account-bot/bin/webex-codex-launcher
StandardInput=socket
StandardOutput=socket
StandardError=journal
Restart=no
TimeoutStartSec=15s
TimeoutStopSec=15s
RuntimeMaxSec=3660s
KillMode=control-group
OOMPolicy=kill
UMask=0077
ReadOnlyPaths=/proc
ReadOnlyPaths=/run/systemd
ReadWritePaths=/var/lib/webex-generic-account-bot/codex-input-staging
ReadWritePaths=/var/lib/webex-codex-runtime-inputs

NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
PrivateDevices=true
PrivateIPC=true
PrivateNetwork=true
ProtectClock=true
ProtectControlGroups=true
ProtectHostname=true
ProtectKernelLogs=true
ProtectKernelModules=true
ProtectKernelTunables=true
ProtectProc=ptraceable
ProcSubset=pid
RestrictAddressFamilies=AF_UNIX
RestrictNamespaces=true
RestrictRealtime=true
RestrictSUIDSGID=true
LockPersonality=true
MemoryDenyWriteExecute=true
CapabilityBoundingSet=CAP_SYS_PTRACE CAP_SETPCAP
AmbientCapabilities=
SystemCallArchitectures=native
IPAddressDeny=any
KeyringMode=private
TasksMax=64
CPUQuota=100%
LimitNOFILE=128
LimitNPROC=64
LimitFSIZE=2112M
MemoryMax=256M
MemorySwapMax=0
