# deploy/systemd/webex-codex-activation-renew.service
[Unit]
Description=Renew Webex Codex activation receipt
Documentation=man:systemd.service(5)
Wants=network-online.target
Requires=webex-codex-launcher.socket webex-config-pull-worker.service
After=network-online.target systemd-tmpfiles-setup.service webex-codex-launcher.socket webex-config-pull-worker.service
Before=webex-generic-account-bot.service
ConditionPathExists=/sys/fs/cgroup/cgroup.controllers

[Service]
Type=oneshot
RemainAfterExit=no
User=root
Group=root
WorkingDirectory=/
ExecStart=/opt/webex-generic-account-bot/bin/webex-codex-activation renew
Restart=no
TimeoutStartSec=5400s
TimeoutStopSec=15s
RuntimeMaxSec=5400s
KillMode=control-group
OOMPolicy=kill
UMask=0077
StandardInput=null
StandardOutput=journal
StandardError=journal
ReadOnlyPaths=/proc
ReadOnlyPaths=/run/systemd
ReadWritePaths=/run/webex-codex-activation
ReadWritePaths=/run/webex-codex-canary
ReadWritePaths=/var/lib/webex-generic-account-bot/codex-input-staging
ReadWritePaths=/var/lib/webex-generic-account-bot/canary-fixtures
ReadWritePaths=/var/lib/webex-codex-runtime-inputs

NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
PrivateDevices=true
PrivateIPC=true
PrivateNetwork=false
ProtectClock=true
ProtectControlGroups=true
ProtectHostname=true
ProtectKernelLogs=true
ProtectKernelModules=true
ProtectKernelTunables=true
ProtectProc=default
ProcSubset=all
RestrictNamespaces=true
RestrictRealtime=true
RestrictSUIDSGID=true
LockPersonality=true
MemoryDenyWriteExecute=true
CapabilityBoundingSet=CAP_CHOWN
AmbientCapabilities=
SystemCallArchitectures=native
KeyringMode=private
TasksMax=128
CPUQuota=200%
LimitNOFILE=256
LimitNPROC=128
LimitFSIZE=64M
MemoryMax=512M
MemorySwapMax=0
