# deploy/systemd/webex-generic-account-bot.service
[Unit]
Description=Webex generic account bot
Documentation=man:systemd.service(5)
Wants=network-online.target
After=network-online.target systemd-tmpfiles-setup.service
ConditionPathExists=/sys/fs/cgroup/cgroup.controllers
ConditionFileIsExecutable=/opt/webex-generic-account-bot/bin/webex-generic-account-bot
ConditionFileNotEmpty=/var/lib/webex-generic-account-bot/rendered/production.toml
ConditionFileNotEmpty=/var/lib/webex-headless-access/access-token

[Service]
Type=simple
User=webex-generic-account-bot
Group=webex-generic-account-bot
SupplementaryGroups=
WorkingDirectory=/var/lib/webex-generic-account-bot/codex-workspace
Environment=PATH=/usr/local/bin:/usr/bin
EnvironmentFile=/etc/webex-generic-account-bot/bot.env
ExecStart=/opt/webex-generic-account-bot/bin/webex-generic-account-bot --config /var/lib/webex-generic-account-bot/rendered/production.toml
Restart=on-failure
RestartSec=5s
TimeoutStartSec=60s
TimeoutStopSec=120s
KillMode=control-group
OOMPolicy=kill
UMask=0077
StandardInput=null
StandardOutput=journal
StandardError=journal
ReadOnlyPaths=/proc
ReadOnlyPaths=/run/systemd
ReadOnlyPaths=/var/lib/webex-generic-account-bot/rendered
ReadOnlyPaths=/var/lib/webex-headless-access
ReadOnlyPaths=/etc/webex-generic-account-bot
InaccessiblePaths=-/run/webex-config-deploy
ReadWritePaths=/var/lib/webex-generic-account-bot/state
ReadWritePaths=/var/lib/webex-generic-account-bot/codex-home
ReadWritePaths=/var/lib/webex-generic-account-bot/codex-workspace

NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
PrivateDevices=true
PrivateIPC=true
RemoveIPC=true
PrivateNetwork=false
ProtectClock=true
ProtectControlGroups=true
ProtectHostname=true
ProtectKernelLogs=true
ProtectKernelModules=true
ProtectKernelTunables=true
ProtectProc=invisible
ProcSubset=all
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
RestrictNamespaces=true
RestrictRealtime=true
RestrictSUIDSGID=true
LockPersonality=true
MemoryDenyWriteExecute=false
CapabilityBoundingSet=
AmbientCapabilities=
SystemCallArchitectures=native
KeyringMode=private
TasksMax=256
CPUQuota=200%
LimitNOFILE=1024
LimitNPROC=256
LimitFSIZE=2304M
LimitCORE=0
MemoryMax=2G
MemorySwapMax=0

[Install]
WantedBy=multi-user.target
