# deploy/systemd/webex-generic-account-bot.service.d/10-codex-launcher.conf
[Service]
SupplementaryGroups=webex-codex-launch
ReadWritePaths=/var/lib/webex-generic-account-bot/codex-input-staging/pending
