# deploy/systemd/webex-generic-account-bot.service.d/10-config-pull.conf
[Service]
SupplementaryGroups=webex-config-pull
